# nodelink-agent.service
[Unit]
Description=Nodelink Agent
Documentation=https://github.com/mooncorn/nodelink
After=network.target
Wants=network.target

[Service]
Type=simple
User=nodelink
Group=nodelink
ExecStart=/usr/local/bin/nodelink-agent -agent_id=${AGENT_ID} -agent_token=${AGENT_TOKEN} -address=${SERVER_ADDRESS}
Restart=always
RestartSec=5
StartLimitInterval=60s
StartLimitBurst=3

# Environment file (will be created by deployment script)
EnvironmentFile=-/etc/nodelink/agent.env

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/log/nodelink /var/lib/nodelink

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=nodelink-agent

[Install]
WantedBy=multi-user.target
