# agent/nodelink-updater.service
[Unit]
Description=Nodelink Agent Updater
Documentation=https://github.com/mooncorn/nodelink
After=network.target
Wants=network.target

[Service]
Type=simple
User=nodelink
Group=nodelink
ExecStart=/usr/local/bin/nodelink-updater -current-version=${AGENT_VERSION} -check-interval=30m -agent-binary=/usr/local/bin/nodelink-agent -repo-owner=mooncorn -repo-name=nodelink
Restart=always
RestartSec=30
StartLimitInterval=300s
StartLimitBurst=5

# Environment file (will be created by deployment script)
EnvironmentFile=-/etc/nodelink/agent.env

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/usr/local/bin /var/log/nodelink /var/lib/nodelink

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=nodelink-updater

[Install]
WantedBy=multi-user.target
